# crontab
*/5 * * * * docker exec shblg_web_1 python3 manage.py fetchmail > /dev/null 2>&1
